# debian/bustd.service
[Unit]
Description=Bustd OOM Daemon

[Service]
ExecStart=/usr/bin/bustd --no-daemon
Restart=on-failure

[Install]
WantedBy=multi-user.target